# clevo-fan-control_0.0.1-1_amd64/etc/xdg/autostart/clevo-fan-control.desktop
[Desktop Entry]
Type=Application
Name=clevo-fan-control
Comment=Driver for fans on Clevo machines
Path=/usr/local/bin
Exec=clevo-fan-control
Terminal=false
